# scripts/steps/rc.local
#!/bin/sh
/etc/rc.d/scripts/apt-cache-on-tmp
